# systemd/aria-intraday.service
[Unit]
Description=ARIA — intraday snapshot
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/Documents/projects/ha-aria
ExecStart=/bin/bash -c '. ~/.env && exec %h/Documents/projects/ha-aria/.venv/bin/aria snapshot-intraday'
StandardOutput=append:%h/.local/log/aria.log
StandardError=append:%h/.local/log/aria.log
MemoryHigh=1G
MemoryMax=2G
